# deploy/timeweb/padel-brazzers-backup.service
[Unit]
Description=Back up Padel Brazzers production store

[Service]
Type=oneshot
ExecStart=/srv/padel-brazzers/config/backup.sh
