# deploy/timeweb/padel-brazzers-backup.timer
[Unit]
Description=Daily Padel Brazzers production store backup

[Timer]
OnCalendar=*-*-* 02:25:00 UTC
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
